# supervisord.conf
[supervisord]
nodaemon=true

[program:ta2]
command=/usr/bin/python /fake-ta2.py
